# admin-tools/crontab
1 0 * * * /home/ubuntu/bot/admin-tools/alerting.sh
